# services/testnet.timer
# testnet.timer

[Unit]
Description=testnet timer

[Timer]
OnUnitActiveSec=100s
OnBootSec=10s

[Install]
WantedBy=timers.target
