# services/testnet.service
# testnet.service

[Unit]
Description=Testnet service daemon with logs for journalctl

[Service]
Type=simple
ExecStart=%h/projects/oss/testnet/testnet

[Install]
WantedBy=default.target
